# matedisplayrule.rules
SUBSYSTEM=="drm", ACTION=="change", RUN+="INSERT PATH TO mpanup.sh"